# scripts/systemd/trading-mike-morning.service
[Unit]
Description=Mike's 9:06 ET morning analysis — invokes Claude Code /mike-morning
Documentation=file:///home/tianyizhang/opus%20trading/.claude/commands/mike-morning.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tianyizhang/opus trading
Environment=PATH=/home/tianyizhang/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/tianyizhang/.local/bin/claude -p "/mike-morning" --dangerously-skip-permissions
TimeoutStartSec=1200
StandardOutput=append:/home/tianyizhang/opus trading/logs/mike-morning.log
StandardError=append:/home/tianyizhang/opus trading/logs/mike-morning.log
