# scripts/systemd/trading-sector-evenings.service
[Unit]
Description=Trading sector-evenings orchestrator — fans out 11 sector-evening attribution reviews
Documentation=file://@REPO_URL_PATH@/scripts/run_evening_orchestrator.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash "@REPO_PATH@/scripts/run_evening_orchestrator.sh"
TimeoutStartSec=2700
StandardOutput=append:@REPO_PATH@/logs/evening-orchestrator.log
StandardError=append:@REPO_PATH@/logs/evening-orchestrator.log
