# scripts/systemd/trading-news-ingest.service
[Unit]
Description=Trading news ingestor — Massive (Benzinga) → Postgres news_items
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
EnvironmentFile=@REPO_PATH@/.env
TimeoutStartSec=300
ExecStart=/bin/bash -c 'cd "@REPO_PATH@" && exec .venv/bin/python scripts/ingest_news.py'
StandardOutput=append:@REPO_PATH@/logs/news-ingest.log
StandardError=append:@REPO_PATH@/logs/news-ingest.log
