# scripts/systemd/trading-weekly-tune.service
[Unit]
Description=Trading weekly model-tune cycle — Python pipeline; tunes 11 sector model portfolios
Documentation=file://@REPO_URL_PATH@/scripts/run_weekly_orchestrator.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash "@REPO_PATH@/scripts/run_weekly_orchestrator.sh"
TimeoutStartSec=10800
StandardOutput=append:@REPO_PATH@/logs/weekly-tune-orchestrator.log
StandardError=append:@REPO_PATH@/logs/weekly-tune-orchestrator.log

[Install]
WantedBy=multi-user.target
