# scripts/systemd/trading-cassidy-evening.service
[Unit]
Description=Cassidy's 11:00 PM MST evening risk review — invokes Claude Code /cassidy-evening
Documentation=file://@REPO_URL_PATH@/.claude/commands/cassidy-evening.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=@CLAUDE_BIN@ -p "/cassidy-evening" --dangerously-skip-permissions
TimeoutStartSec=1200
StandardOutput=append:@REPO_PATH@/logs/cassidy-evening.log
StandardError=append:@REPO_PATH@/logs/cassidy-evening.log
