# scripts/systemd/trading-concierge.user.service
[Unit]
Description=Trading Concierge (local-LLM-backed Telegram bot, user mode)
Documentation=file://@REPO_URL_PATH@/concierge/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO_PATH@
EnvironmentFile=@REPO_PATH@/.env
# systemd splits ExecStart on whitespace even when the binary is quoted —
# wrap in bash -c so the space in "@REPO_PATH@" survives.
ExecStart=/bin/bash -c 'cd "@REPO_PATH@" && exec .venv/bin/python -m concierge.service'
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=10

[Install]
# Mode C analog: runs while the user has an active session. To survive
# logout/reboot without re-login, run:  loginctl enable-linger "$USER"
WantedBy=default.target
