# scripts/systemd/trading-forecast-resolver.service
[Unit]
Description=Trading forecast outcome resolver — fills realized_return_pct on horizon-elapsed forecasts
Documentation=file:///home/tianyizhang/opus%%20trading/scripts/run_forecast_resolver.py
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tianyizhang/opus trading
Environment=PATH=/home/tianyizhang/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -c 'cd "/home/tianyizhang/opus trading" && exec .venv/bin/python scripts/run_forecast_resolver.py'
TimeoutStartSec=600
StandardOutput=append:/home/tianyizhang/opus trading/logs/forecast-resolver.log
StandardError=append:/home/tianyizhang/opus trading/logs/forecast-resolver.log

[Install]
WantedBy=multi-user.target
