# scripts/systemd/trading-sector-evenings.service
[Unit]
Description=Trading sector-evenings orchestrator — fans out 11 sector-evening attribution reviews
Documentation=file:///home/tianyizhang/opus%20trading/scripts/run_evening_orchestrator.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tianyizhang/opus trading
Environment=PATH=/home/tianyizhang/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash "/home/tianyizhang/opus trading/scripts/run_evening_orchestrator.sh"
TimeoutStartSec=2700
StandardOutput=append:/home/tianyizhang/opus trading/logs/evening-orchestrator.log
StandardError=append:/home/tianyizhang/opus trading/logs/evening-orchestrator.log
