# scripts/systemd/trading-vllm.service
[Unit]
Description=Trading vLLM — local OpenAI-compatible LLM server (Qwen3.6-35B-A3B-FP8)
Documentation=https://docs.vllm.ai/en/latest/serving/openai_compatible_server.html
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
WorkingDirectory=/home/tianyizhang/opus trading
EnvironmentFile=/home/tianyizhang/opus trading/.env
Environment=HF_HUB_DISABLE_TELEMETRY=1
Environment=HF_HUB_DISABLE_PROGRESS_BARS=1
Environment=VLLM_LOGGING_LEVEL=INFO
# 10-minute start budget — first launch may load 38 GB of FP8 weights into VRAM.
TimeoutStartSec=600
## Requires `python3.12-devel` so vLLM's Inductor/Triton JIT can compile FP8
## kernels. If that package is ever uninstalled, add `--enforce-eager` to the
## ExecStart below (vLLM still serves, just ~30% slower per token).
## Aliases let the `claude` CLI hit /v1/messages with its native model name
## (claude-opus-4-7 by default; sometimes claude-sonnet-* on subagents) and
## get routed to the local Qwen backend without a separate proxy. vLLM 0.20+
## also natively supports the Anthropic /v1/messages shape, so no LiteLLM
## proxy needed.
## Atlas-review and other hourly skills routinely build ~45-50K token prompts
## (system prompt + skill markdown + journal history + tool result interpolations).
## Qwen3-32B's native context is 40,960 tokens; YARN extends it to 131,072.
## We cap at 80K so KV-cache stays under ~25GB even with 4 concurrent sector agents.
ExecStart=/bin/bash -c 'cd "/home/tianyizhang/opus trading" && exec .venv-vllm/bin/vllm serve "${LOCAL_MODEL}" \
    --served-model-name "${LOCAL_MODEL}" claude-opus-4-7 claude-sonnet-4-5-20250929 claude-sonnet-4-6 claude-haiku-4-5-20251001 claude-3-5-sonnet-20241022 \
    --quantization fp8 \
    --gpu-memory-utilization 0.92 \
    --max-model-len 81920 \
    --hf-overrides '"'"'{"rope_scaling":{"rope_type":"yarn","factor":2.0,"original_max_position_embeddings":40960}}'"'"' \
    --enable-auto-tool-choice \
    --tool-call-parser hermes \
    --host 127.0.0.1 \
    --port 8000'
Restart=on-failure
RestartSec=15
# A vLLM crash-loop on missing weights or OOM should not retry forever.
StartLimitIntervalSec=600
StartLimitBurst=4
StandardOutput=append:/home/tianyizhang/opus trading/logs/vllm.log
StandardError=append:/home/tianyizhang/opus trading/logs/vllm.log

[Install]
WantedBy=default.target
