# scripts/systemd/trading-hourly-review.service
[Unit]
Description=Trading hourly orchestrator — fans out 10 sector reviews + mike-allocator, then /hourly-review heartbeat
Documentation=file://@REPO_URL_PATH@/scripts/run_hourly_orchestrator.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash "@REPO_PATH@/scripts/run_hourly_orchestrator.sh"
TimeoutStartSec=2700
StandardOutput=append:@REPO_PATH@/logs/hourly-orchestrator.log
StandardError=append:@REPO_PATH@/logs/hourly-orchestrator.log
