# scripts/systemd/trading-news-ingest.service
[Unit]
Description=Trading news ingestor — Massive (Benzinga) → Postgres news_items
Documentation=file:///home/tianyizhang/opus%20trading/scripts/ingest_news.py
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tianyizhang/opus trading
EnvironmentFile=/home/tianyizhang/opus trading/.env
TimeoutStartSec=300
ExecStart=/bin/bash -c 'cd "/home/tianyizhang/opus trading" && exec .venv/bin/python scripts/ingest_news.py'
StandardOutput=append:/home/tianyizhang/opus trading/logs/news-ingest.log
StandardError=append:/home/tianyizhang/opus trading/logs/news-ingest.log
