# scripts/systemd/trading-refresh-agent-state.service
[Unit]
Description=Trading per-agent state hourly snapshot — deterministic, no LLM
Documentation=file://@REPO_URL_PATH@/scripts/refresh_agent_state.sh
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash "@REPO_PATH@/scripts/refresh_agent_state.sh"
TimeoutStartSec=600
StandardOutput=append:@REPO_PATH@/logs/refresh_agent_state.log
StandardError=append:@REPO_PATH@/logs/refresh_agent_state.log
