# scripts/systemd/trading-embed-sweeper.service
[Unit]
Description=Trading embed sweeper — backfill NULL-embedding news_items rows
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
EnvironmentFile=@REPO_PATH@/.env
TimeoutStartSec=300
ExecStart=/bin/bash -c 'cd "@REPO_PATH@" && exec .venv/bin/python scripts/embed_news_sweeper.py'
StandardOutput=append:@REPO_PATH@/logs/embed-sweeper.log
StandardError=append:@REPO_PATH@/logs/embed-sweeper.log
