# scripts/systemd/trading-concierge.service
[Unit]
Description=Trading Concierge (local-LLM-backed Telegram bot, user mode)
Documentation=file:///home/tianyizhang/opus%20trading/concierge/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tianyizhang/opus trading
EnvironmentFile=/home/tianyizhang/opus trading/.env
# systemd splits ExecStart on whitespace even when the binary is quoted —
# wrap in bash -c so the space in "/home/tianyizhang/opus trading" survives.
ExecStart=/bin/bash -c 'cd "/home/tianyizhang/opus trading" && exec .venv/bin/python -m concierge.service'
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=10

[Install]
# Mode C analog: runs while the user has an active session. To survive
# logout/reboot without re-login, run:  loginctl enable-linger tianyizhang
WantedBy=default.target
