# scripts/systemd/trading-refresh-agent-state.timer
[Unit]
Description=Trading per-agent state hourly snapshot timer — every hour at :05

[Timer]
# Replaces the cron entry in install_scheduled_tasks.sh, which silently failed
# whenever REPO_PATH contained a space (cron tried to exec the first
# whitespace-delimited token). systemd handles spaces in paths correctly.
OnCalendar=*-*-* *:05:00
Persistent=true
AccuracySec=30s
Unit=trading-refresh-agent-state.service

[Install]
WantedBy=timers.target
