# scripts/systemd/trading-mike-midday.service
[Unit]
Description=Mike's 11:00 ET midday review — invokes Claude Code /mike-midday
Documentation=file://@REPO_URL_PATH@/.claude/commands/mike-midday.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=@CLAUDE_BIN@ -p "/mike-midday" --dangerously-skip-permissions
TimeoutStartSec=1200
StandardOutput=append:@REPO_PATH@/logs/mike-midday.log
StandardError=append:@REPO_PATH@/logs/mike-midday.log
