# scripts/systemd/trading-news-ingest.timer
[Unit]
Description=Trading news ingest timer — every 15 min on weekdays (script no-ops cheaply outside RTH)

[Timer]
# Every 15 min on the quarter-hour, weekdays only. systemd's calendar specifier
# rejects step ranges across hour boundaries (`HH..HH/0:15:00` is a bad-setting),
# so we keep it simple and let the script handle market-hours filtering itself.
OnCalendar=Mon..Fri *-*-* *:00,15,30,45:00
Persistent=true
AccuracySec=30s
Unit=trading-news-ingest.service

[Install]
WantedBy=timers.target
