# scripts/systemd/trading-news-ingest.timer
[Unit]
Description=Trading news ingest timer — every 10 min on weekdays 24h (covers pre/post-market earnings)

[Timer]
# Every 10 min, weekdays 24h. Pre-market and after-hours news (earnings releases,
# foreign-market opens, late M&A drops) move next-day positions — cassidy-evening
# at 23:00 MT and mike-morning at 09:06 ET both want fresh items.
OnCalendar=Mon..Fri *-*-* *:00/10:00
Persistent=true
AccuracySec=30s
Unit=trading-news-ingest.service

[Install]
WantedBy=timers.target
