# scripts/systemd/trading-embed-sweeper.service
[Unit]
Description=Trading embed sweeper — backfill NULL-embedding news_items rows
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/tianyizhang/opus trading
EnvironmentFile=/home/tianyizhang/opus trading/.env
TimeoutStartSec=300
ExecStart=/bin/bash -c 'cd "/home/tianyizhang/opus trading" && exec .venv/bin/python scripts/embed_news_sweeper.py'
StandardOutput=append:/home/tianyizhang/opus trading/logs/embed-sweeper.log
StandardError=append:/home/tianyizhang/opus trading/logs/embed-sweeper.log
