# scripts/systemd/trading-mike-morning.service
[Unit]
Description=Mike's 9:06 ET morning analysis — local-LLM harness for /mike-morning
Documentation=file:///home/tianyizhang/opus%20trading/.claude/commands/mike-morning.md
After=network-online.target trading-vllm.service trading-llm-proxy.service
Wants=network-online.target trading-vllm.service trading-llm-proxy.service

[Service]
Type=oneshot
WorkingDirectory=/home/tianyizhang/opus trading
EnvironmentFile=/home/tianyizhang/opus trading/.env
Environment=PATH=/home/tianyizhang/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/bin/bash -c 'cd "/home/tianyizhang/opus trading" && exec .venv/bin/python scripts/run_harness_local.py mike-morning'
TimeoutStartSec=1800
StandardOutput=append:/home/tianyizhang/opus trading/logs/mike-morning.log
StandardError=append:/home/tianyizhang/opus trading/logs/mike-morning.log
