# scripts/systemd/trading-mike-morning.service
[Unit]
Description=Mike's 9:06 ET morning analysis — invokes Claude Code /mike-morning
Documentation=file://@REPO_URL_PATH@/.claude/commands/mike-morning.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=@REPO_PATH@
Environment=PATH=@CLAUDE_BIN_DIR@:/usr/local/bin:/usr/bin:/bin
ExecStart=@CLAUDE_BIN@ -p "/mike-morning" --dangerously-skip-permissions
TimeoutStartSec=1200
StandardOutput=append:@REPO_PATH@/logs/mike-morning.log
StandardError=append:@REPO_PATH@/logs/mike-morning.log
