# scripts/systemd/trading-desk.target
[Unit]
Description=Trading desk — all long-running daemons
Documentation=https://github.com/LofiChptr6/ProjectCorvus

# The daemons in Wants= are the long-running components that make up "the desk
# is up". Timer-triggered oneshots (hourly orchestrator, cassidy evening, mike
# morning/midday, refresh_agent_state, news ingest, sector evenings, embed
# sweeper, weekly tune) are NOT here — they fire independently and don't need
# to be co-managed with the daemons. Stopping this target only stops the
# daemons; timers continue to fire and try to run (failing fast against missing
# IBKR / vLLM is fine).
Wants=ibkr-daemon.service
Wants=trading-vllm.service
Wants=trading-vllm-embed.service
Wants=trading-llm-proxy.service
Wants=trading-concierge.service
Wants=trading-dashboard.service

# Order: nothing here; each daemon manages its own After= chain. The target is
# purely a grouping handle for "start/stop the desk."
After=network-online.target

[Install]
WantedBy=default.target
